# files/system/usr/lib/systemd/system/flatpak-system-update.timer.d/60-update-schedule.conf
[Unit]
Description=Flatpak Update Trigger

[Timer]
OnCalendar=
OnCalendar=daily
OnBootSec=10m
RandomizedDelaySec=30m
