# files/system/usr/lib/systemd/system/flatpak-system-update.timer.d/override.conf
[Unit]
Description=Flatpak Update Trigger (Every 15 mins)

[Timer]
OnCalendar=
OnCalendar=*:0/15
OnBootSec=5m
RandomizedDelaySec=5m