# files/system/usr/lib/systemd/system/flatpak-system-update.timer.d/override.conf
[Unit]
Description=Flatpak Update Trigger  (Hourly)

[Timer]
OnCalendar=
OnCalendar=*:0/30
# OnCalendar=hourly
OnBootSec=10m
RandomizedDelaySec=5m