# files/system/etc/systemd/system/flatpak-system-update.timer.d/99-custom-schedule.conf
[Unit]
Description=Flatpak Update Trigger

[Timer]
OnCalendar=
OnCalendar=*:0/30
OnBootSec=10m
RandomizedDelaySec=5m