# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.timer.d/60-custom-schedule.conf
[Unit]
Description=rpm-ostree Automatic Update Trigger

[Timer]
OnCalendar=
OnCalendar=hourly
OnBootSec=10m
RandomizedDelaySec=5m