# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.timer.d/60-update-schedule.conf
[Unit]
Description=rpm-ostree Automatic Update Trigger

[Timer]
OnCalendar=
OnCalendar=daily
OnBootSec=10m
RandomizedDelaySec=30m
